# docker/queue-worker.ini
[program:clubhouse-worker]
process_name=%(program_name)s_%(process_num)02d
command=sh /usr/bin/clubhouse-worker
autostart=true
autorestart=true
numprocs=1
startretries=5
stdout_events_enabled=1
redirect_stderr=true
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:cron]
process_name=%(program_name)s_%(process_num)02d
command=sh /usr/bin/clubhouse-scheduler
autostart=true
autorestart=true
user=www-data
numprocs=1
startretries=5
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
