# docker/queue-worker.ini
[program:clubhouse-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/application/artisan queue:work --sleep=60 --tries=3 --memory=1024
autostart=true
autorestart=true
numprocs=1
startretries=10
stdout_events_enabled=1
redirect_stderr=true
user=www-data

[program:schedule-worker]
process_name=%(program_name)s_%(process_num)02d
command=sh /usr/bin/schedule-worker
autostart=true
autorestart=true
numprocs=1
startretries=5
stdout_events_enabled=1
redirect_stderr=true
user=www-data
