# docker/supervisord-nginx.ini
[program:nginx]
command=start-nginx -g "daemon off;"
priority=200
autorestart=true
