# docker/supervisord-php-fpm.ini
[program:php-fpm]
command=php-fpm --nodaemonize
autorestart=true
autostart=false
