# docker/supervisord-nginx.ini
[program:nginx]
command=start-nginx -g "daemon off;"
priority=200
autorestart=true

[program:redis]
command=/usr/bin/redis-server
user=root
autostart=false
stopsignal=QUIT
