# docker/supervisord-php-octane.ini
[program:php-octane]
directory=/var/www/application
command=php artisan octane:start
autostart=false
autorestart=true
redirect_stderr=true
priority=100
